# systemd/tas2781-bypass-resume.service
[Unit]
Description=Reapply TAS2781 speaker bypass after resume
Before=sleep.target
StopWhenUnneeded=yes
ConditionPathExists=/usr/local/lib/tas2781-bypass/tas2781-bypass-apply.sh

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/true
ExecStop=/usr/local/lib/tas2781-bypass/tas2781-bypass-apply.sh --resume
TimeoutStopSec=30
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=sleep.target
